# dist/linux/systemd/joysockmapper.socket
[Unit]
Description=JoyShockMapper Config Reload Socket
Documentation=https://github.com/aapokaapo/JoyShockMapper
PartOf=joysockmapper.service

[Socket]
ListenStream=/run/joysockmapper.sock
Accept=no

[Install]
WantedBy=sockets.target
