# dist/linux/systemd/joyshockmapper.socket
[Unit]
Description=JoyShockMapper Config Reload Socket
Documentation=https://github.com/aapokaapo/JoyShockMapper
PartOf=joyshockmapper.service

[Socket]
ListenStream=%t/joyshockmapper.sock
Accept=no

[Install]
WantedBy=sockets.target
